# Ads-manager/deploy/systemd/meta-ads-worker.service
[Unit]
Description=Meta Ads Dev Worker
After=network.target postgresql.service meta-ads-api.service
Wants=postgresql.service

[Service]
Type=simple
User=metaads
Group=metaads
WorkingDirectory=/opt/meta-ads-dev/app
EnvironmentFile=/opt/meta-ads-dev/app/.env
ExecStart=/usr/bin/node dist/worker.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/meta-ads-dev/app
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
